# bad-parents.service
[Unit]
Description=bad_parents — TikTok block controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=simon
WorkingDirectory=/home/simon/projects/bad_parents
ExecStart=/usr/bin/python3 /home/simon/projects/bad_parents/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
